# SQL_Parser/conf/sql-parser-supervisor.ini
[program:sql-parser-uwsgi]
command = uwsgi /home/Sql_Parser/conf/sql-parser-uwsgi.ini
autostart=true
